# collector/market-data-fundamentals.timer
[Unit]
Description=Weekly Fundamental Data Collection

[Timer]
OnCalendar=Sat 08:00
Persistent=true

[Install]
WantedBy=timers.target
